# systemd/skyline-server.service
[Unit]
Description=Skyline Web UI & Profile Server
Documentation=https://skyline.projex.cc/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/skyline-server --listen=localhost:19190 --storage=%h/.skyline/profiles.enc.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-%h/.skyline/skyline.env

# Security
NoNewPrivileges=true
PrivateTmp=true

# Working directory
WorkingDirectory=%h/.skyline

[Install]
WantedBy=default.target
