# systemd/skyline.service
[Unit]
Description=Skyline MCP Server
Documentation=https://skyline.projex.cc/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/skyline --config=%h/.skyline/config.yaml --transport=http --listen=:8191
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Security
NoNewPrivileges=true
PrivateTmp=true

# Working directory
WorkingDirectory=%h/.skyline

[Install]
WantedBy=default.target
